# crontab
*/5 * * * * curl -X POST http://xxx.flycast/start-task
